# etc/profile.d/20-bat-manpager.sh
export MANPAGER="sh -c 'col -bx | LESS="RKC" bat -l man --decorations=never --paging=always'"
export MANROFFOPT="-c"
